# lib/ipkg/info/php7-fastcgi.list
/
/etc/
/etc/config/
/etc/config/php7-fastcgi
/etc/init.d/
/etc/init.d/php7-fastcgi
